# roles/instafap/files/upstart/instafap.conf
start on runlevel [2345]
stop on runlevel [!2345]

respawn
respawn limit 10 5

setuid root
chdir /opt/instafap

env NODE_ENV="production"
env INSTAFAP_SECRET="****************"
env AWS_ACCESS_KEY="****************"
env AWS_SECRET_KEY="****************"
env FACEBOOK_CLIENT_ID="****************"
env FACEBOOK_CLIENT_SECRET="****************"

exec npm start